# contrib/init/dogecoind.conf
description "Whippet Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env WHTCOIND_BIN="/usr/bin/whippetd"
env WHTCOIND_USER="whippet"
env WHTCOIND_GROUP="whippet"
env WHTCOIND_PIDDIR="/var/run/whippetd"
# upstart can't handle variables constructed with other variables
env WHTCOIND_PIDFILE="/var/run/whippetd/whippetd.pid"
env WHTCOIND_CONFIGFILE="/etc/whippet/whippet.conf"
env WHTCOIND_DATADIR="/var/lib/whippetd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # whippetd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$WHTCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run whippetd."
        echo "The setting must appear in $WHTCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Whippet Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$WHTCOIND_PIDDIR"
    chmod 0755 "$WHTCOIND_PIDDIR"
    chown $WHTCOIND_USER:$WHTCOIND_GROUP "$WHTCOIND_PIDDIR"
    chown $WHTCOIND_USER:$WHTCOIND_GROUP "$WHTCOIND_CONFIGFILE"
    chmod 0660 "$WHTCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$WHTCOIND_PIDFILE" \
    --chuid $WHTCOIND_USER:$WHTCOIND_GROUP \
    --exec "$WHTCOIND_BIN" \
    -- \
    -pid="$WHTCOIND_PIDFILE" \
    -conf="$WHTCOIND_CONFIGFILE" \
    -datadir="$WHTCOIND_DATADIR" \
    -disablewallet \
    -daemon

